# daemons/lprivileges.path
# version 0.3
[Unit]
Description=Monitor Linux Privileges Request file

[Path]
PathModified=/var/tmp/lprivileges/lprivileges
Unit=lprivileges.service

[Install]
WantedBy=multi-user.target
