# daemons/lprivileges.service
# version 0.3
[Unit]
Description=Executes lprivilegestriage when a file request si made.

[Service]
ExecStart=/usr/local/bin/lprivilegestriage
